# deploy/systemd/ai-soc-qdrant-detection-case-backfill.timer
[Unit]
Description=Run AI SOC Qdrant detection/case memory backfill daily

[Timer]
OnCalendar=*-*-* 03:35:00
AccuracySec=5min
Persistent=true

[Install]
WantedBy=timers.target
