# deploy/systemd/ai-soc-qdrant-historical-backfill.service
[Unit]
Description=AI SOC Qdrant historical incident memory backfill
After=ai-soc-api.service
Wants=ai-soc-api.service

[Service]
Type=oneshot
WorkingDirectory=/home/lele/lab/ai-soc-assistant
Environment=PYTHONPATH=/home/lele/lab/ai-soc-assistant
ExecStart=/home/lele/lab/ai-soc-assistant/.venv/bin/python /home/lele/lab/ai-soc-assistant/scripts/index_historical_incidents_to_qdrant.py --apply --limit 10000
User=lele
