# systemd/ai-soc-api.service
[Unit]
Description=Sovereign AI SOC FastAPI Backend
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=lele
Group=docker
WorkingDirectory=/home/lele/lab/ai-soc-assistant
Environment="PYTHONUNBUFFERED=1"
ExecStart=/home/lele/lab/ai-soc-assistant/.venv/bin/uvicorn api:app --host 0.0.0.0 --port 8008
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target

